# sources/meta-digi/meta-digi-arm/dynamic-layers/freescale-layer/recipes-kernel/kernel-modules/kernel-module-nxp-wlan/watch-regdomain.timer
[Unit]
Description=Check regdomain every 5 seconds

[Timer]
OnBootSec=5sec
OnUnitInactiveSec=5sec
AccuracySec=1s
Unit=watch-regdomain.service

[Install]
WantedBy=timers.target
