# sources/meta-digi/meta-digi-arm/dynamic-layers/freescale-layer/recipes-kernel/kernel-modules/kernel-module-nxp-wlan/watch-regdomain.service
[Unit]
Description=Monitor for global regulatory domain changes
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/watch_regdomain.sh
